# openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/p4nd4pr0t0c01/Projects/openclaw/scripts/run-node.mjs gateway
WorkingDirectory=/home/p4nd4pr0t0c01/Projects/openclaw
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PATH=/home/p4nd4pr0t0c01/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
StandardOutput=append:/tmp/openclaw-service.log
StandardError=append:/tmp/openclaw-service.log

[Install]
WantedBy=default.target
